# GdoSnifferKevelaer.timer
[Unit]
Description=GdoSniffer service for Kevelaer

[Timer]
OnCalendar=*-*-* 00/2:0:0

[Install]
WantedBy=basic.target
