# systemd/gdo-sniffer-kevelaer.timer
[Unit]
Description=GdoSniffer service for Kevelaer

[Timer]
OnCalendar=*-*-* 00/2:00:00

[Install]
WantedBy=basic.target
