# systemd/gdo-sniffer.timer
[Unit]
Description=GdoChecker service

[Timer]
OnCalendar=*-*-* 8..22/2:00:00

[Install]
WantedBy=basic.target
